# bamboo-python-lvgl.service
[Unit]
Description=Bamboo Python LVGL Detection System
After=network.target graphical-session.target
Wants=network.target

[Service]
Type=simple
User=jetson
Group=jetson
WorkingDirectory=/opt/bamboo-cut
ExecStart=/opt/bamboo-cut/venv/bin/python python_core/ai_bamboo_system.py
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

# 环境变量设置
Environment="PYTHONPATH=/opt/bamboo-cut:/opt/bamboo-cut/python_core"
Environment="LD_LIBRARY_PATH=/usr/local/cuda/lib64:/opt/bamboo-cut/lib"
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"

# 系统资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bamboo-python-lvgl

[Install]
WantedBy=multi-user.target