# package/DEBIAN/prerm
#!/bin/bash
# 智能切竹机系统卸载前清理脚本

set -e

# 配置变量
SERVICE_USER="bamboo-cut"

echo "正在停止智能切竹机系统服务..."

# 停止所有相关服务
services=("bamboo-cut" "bamboo-cut-kiosk")
for service in "${services[@]}"; do
    if systemctl is-active --quiet "$service.service" 2>/dev/null; then
        echo "停止服务: $service"
        systemctl stop "$service.service" || true
    fi
    
    if systemctl is-enabled --quiet "$service.service" 2>/dev/null; then
        echo "禁用服务: $service"
        systemctl disable "$service.service" || true
    fi
done

# 重新加载 systemd
systemctl daemon-reload || true

# 如果是完全卸载，询问是否删除用户数据
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    echo "准备卸载智能切竹机系统..."
    
    # 停止所有由该用户运行的进程
    if id "$SERVICE_USER" &>/dev/null; then
        pkill -u "$SERVICE_USER" || true
        sleep 2
        pkill -9 -u "$SERVICE_USER" || true
    fi
fi

echo "智能切竹机系统服务已停止"

exit 0 