# deploy/systemd/bamboo-qt-ui.service
[Unit]
Description=Bamboo Qt DeepStream UI (offscreen)
After=multi-user.target network-online.target nvargus-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=root
RuntimeDirectory=bamboo-qt
RuntimeDirectoryMode=0700
WorkingDirectory=/opt/bamboo-qt
Environment=QT_QPA_PLATFORM=offscreen
Environment=QT_QPA_EGLFS_INTEGRATION=
Environment=QT_QPA_EGLFS_KMS_CONFIG=
Environment=GBM_DRIVERS_PATH=
Environment=GBM_BACKEND=
Environment=EGL_PLATFORM=
Environment=XDG_RUNTIME_DIR=/run/bamboo-qt
Environment=QSG_RHI_BACKEND=software
Environment=QT_QUICK_BACKEND=software
Environment=QML2_IMPORT_PATH=/usr/lib/aarch64-linux-gnu/qt6/qml
Environment=QT_PLUGIN_PATH=/usr/lib/aarch64-linux-gnu/qt6/plugins
Environment=GST_DEBUG=2
ExecStartPre=/bin/mkdir -p -m 700 /run/bamboo-qt
ExecStart=/opt/bamboo-qt/bin/bamboo_qt_ui
TimeoutStartSec=30
Restart=always
RestartSec=3
SupplementaryGroups=video render
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
