# deploy/systemd/bamboo-python-gtk4.service
[Unit]
Description=Bamboo Python GTK4 Detection System
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bamboo-cut
Environment="PYTHONPATH=/opt/bamboo-cut"
Environment="PATH=/opt/bamboo-cut/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/opt/bamboo-cut/venv"
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="HOME=/root"
Environment="USER=root"

# 智能环境变量 - 支持多种显示后端
Environment="GDK_BACKEND=wayland,x11"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="DISPLAY=:0"
Environment="XDG_SESSION_TYPE=wayland"

# 直接执行Python应用（智能检测环境）
ExecStart=/opt/bamboo-cut/venv/bin/python /opt/bamboo-cut/python_core/ai_bamboo_system.py

Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target