# deploy/systemd/bamboo-vision.service
[Unit]
Description=Bamboo Vision (jetson-inference + Modbus + HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bamboo-vision
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/bamboo-vision/jetson-inference/python:/opt/bamboo-vision/jetson-inference/build/aarch64/python:$PYTHONPATH
Environment=LD_LIBRARY_PATH=/opt/bamboo-vision/jetson-inference/build/aarch64/lib:/usr/lib/aarch64-linux-gnu/tegra:/usr/lib/aarch64-linux-gnu:$LD_LIBRARY_PATH
ExecStart=/usr/bin/python3 -m bamboo_vision.app --config /opt/bamboo-vision/config/runtime.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
