# deploy/systemd/bamboo-python-gtk4.service
[Unit]
Description=Bamboo Python GTK4 Detection System
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bamboo-cut
Environment="PYTHONPATH=/opt/bamboo-cut"
Environment="PATH=/opt/bamboo-cut/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/opt/bamboo-cut/venv"

# 修复权限问题 - 使用root的运行目录
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="HOME=/root"
Environment="USER=root"

# 禁用dconf以避免权限问题
Environment="DCONF_PROFILE="
Environment="GSETTINGS_BACKEND=memory"

# GTK4环境变量 - 无头模式友好
Environment="GDK_BACKEND=broadway,wayland,x11"
Environment="BROADWAY_DISPLAY=:1"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="DISPLAY=:0"
Environment="XDG_SESSION_TYPE=wayland"

# 直接执行Python应用（智能检测环境）
ExecStart=/opt/bamboo-cut/venv/bin/python /opt/bamboo-cut/python_core/ai_bamboo_system.py

Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target