# deploy/systemd/bamboo-cpp-lvgl.service.in
[Unit]
Description=AI Bamboo Recognition System - C++ LVGL Edition
Documentation=https://github.com/bamboo-industries/bamboo-recognition
After=network.target multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/bin
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/bamboo_integrated --verbose --config @CMAKE_INSTALL_PREFIX@/share/bamboo/config/system_config.yaml
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib:/usr/local/cuda/lib64:/usr/lib/aarch64-linux-gnu"
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="HOME=/root"

# DRM/KMS显示设置
Environment="QT_QPA_PLATFORM=eglfs"
Environment="QT_QPA_EGLFS_INTEGRATION=eglfs_kms"
Environment="QT_QPA_EGLFS_KMS_CONFIG=@CMAKE_INSTALL_PREFIX@/share/bamboo/config/kms.conf"

# 设备权限组
SupplementaryGroups=input video render tty dialout

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
OOMScoreAdjust=-100

# 重启策略
Restart=always
RestartSec=15
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

# 启动前检查和清理
ExecStartPre=/bin/sh -c 'pkill -f bamboo_cu || true'
ExecStartPre=/bin/sh -c 'pkill -f bamboo_integrated || true'
ExecStartPre=/bin/sh -c 'fuser -k /dev/video0 2>/dev/null || true'
ExecStartPre=/bin/sleep 2

# 安全设置（允许访问硬件设备）
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# 日志输出
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target