# deploy/systemd/bamboo-cpp-lvgl.service.in
[Unit]
Description=AI Bamboo Recognition System - C++ LVGL Wayland Edition
Documentation=https://github.com/bamboo-industries/bamboo-recognition
After=network.target multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/bin
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/bamboo_integrated --verbose --config @CMAKE_INSTALL_PREFIX@/share/bamboo/config/system_config.yaml
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib:/usr/local/cuda/lib64:/usr/lib/aarch64-linux-gnu"
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="HOME=/root"
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/etc/bamboo/bamboo-env

# Wayland 显示设置（自动适配）
# Weston 13 使用 wayland-0，Weston 12 使用 wayland-1
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="WAYLAND_DEBUG=0"
Environment="EGL_PLATFORM=wayland"
Environment="QT_QPA_PLATFORM=wayland"
Environment="GDK_BACKEND=wayland"
Environment="SDL_VIDEODRIVER=wayland"
Environment="XDG_SESSION_TYPE=wayland"
Environment="LIBINPUT_DEFAULT_TOUCH_ENABLED=1"

# 设备权限组
SupplementaryGroups=input video render tty dialout

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
OOMScoreAdjust=-100

# 启动前检查和清理
ExecStartPre=/bin/bash -c 'systemctl is-active --quiet bamboo-cpp-lvgl.service && systemctl stop bamboo-cpp-lvgl.service || true'
ExecStartPre=/bin/bash -c 'pkill -9 bamboo_integrated || true'
ExecStartPre=/bin/bash -c 'pkill -9 bamboo_cu || true'
ExecStartPre=/bin/bash -c 'for p in $(lsof -t /dev/video0 2>/dev/null); do kill -9 $p; done || true'
ExecStartPre=/bin/bash -c 'for p in $(lsof -t /dev/video1 2>/dev/null); do kill -9 $p; done || true'
ExecStartPre=/bin/bash -c 'systemctl restart nvargus-daemon'
# 检查并启动 Wayland 合成器（优先级: Weston12 > Sway > Mutter）
ExecStartPre=/bin/bash -c 'if systemctl list-unit-files | grep -q weston12.service; then systemctl is-active --quiet weston12.service || systemctl start weston12.service; elif systemctl list-unit-files | grep -q sway-wayland.service; then systemctl is-active --quiet sway-wayland.service || systemctl start sway-wayland.service; elif systemctl list-unit-files | grep -q mutter-wayland.service; then systemctl is-active --quiet mutter-wayland.service || systemctl start mutter-wayland.service; elif systemctl list-unit-files | grep -q weston.service; then systemctl is-active --quiet weston.service || systemctl start weston.service; fi'
ExecStartPre=/bin/bash -c 'if [ ! -S /run/user/0/wayland-0 ] && [ ! -S /run/user/0/wayland-1 ]; then echo "启动 Wayland 合成器..."; if systemctl list-unit-files | grep -q weston12.service; then systemctl start weston12.service 2>/dev/null || true; elif [ -f /usr/bin/sway ]; then systemctl start sway-wayland.service 2>/dev/null || true; elif [ -f /usr/bin/mutter ]; then systemctl start mutter-wayland.service 2>/dev/null || true; elif [ -f /usr/bin/weston ]; then systemctl start weston.service 2>/dev/null || true; fi; sleep 3; fi'
ExecStartPre=/bin/bash -c 'timeout 30 bash -c "while [ ! -S /run/user/0/wayland-0 ] && [ ! -S /run/user/0/wayland-1 ]; do sleep 1; done" || echo "Wayland socket 不可用，继续使用 fallback 模式"'
ExecStartPre=/bin/sleep 3

# 改进停止命令
ExecStop=/bin/bash -c 'kill -TERM $MAINPID; sleep 2; pkill -9 -g $MAINPID bamboo_integrated || true'

# 更激进的清理
TimeoutStopSec=10
KillMode=control-group
SendSIGKILL=yes
# 重启策略
Restart=always
RestartSec=15
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=mixed

# 安全设置（允许访问硬件设备）
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# 日志输出
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
