# deploy/systemd/bamboo-rtsp-hls.service
[Unit]
Description=Repackage Bamboo RTSP to HLS (LAN playback)
After=bamboo-qt-ui.service network-online.target
Wants=bamboo-qt-ui.service

[Service]
Type=simple
User=root
# Create output dir for HLS segments/playlist
ExecStartPre=/bin/mkdir -p /opt/bamboo-qt/www
# Pull RTSP (local) and generate short HLS segments for low latency
ExecStart=/opt/bamboo-qt/bin/hls-repack.sh
Environment=PIPELINE_ID=rtsp-hls
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
