# deploy/systemd/bamboo-cpp-lvgl.service.in
[Unit]
Description=AI Bamboo Recognition System - C++ LVGL DRM/GBM Edition
Documentation=https://github.com/bamboo-industries/bamboo-recognition
After=network.target multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@CMAKE_INSTALL_PREFIX@/bin
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/bamboo_integrated --verbose --config @CMAKE_INSTALL_PREFIX@/share/bamboo/config/system_config.yaml
ExecReload=/bin/kill -USR1 $MAINPID

# 基础环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib:/usr/local/cuda/lib64:/usr/lib/aarch64-linux-gnu"
Environment="XDG_RUNTIME_DIR=/tmp"
Environment="HOME=/root"
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/etc/bamboo/bamboo-env

# DRM/GBM 设置：禁用 Wayland，直接使用 DRM
Environment="EGL_PLATFORM=drm"
# 可选：强制指定 DRM 设备，留空则由程序内部选择
# Environment="EGL_DRM_DEVICE_FILE=/dev/dri/card1"
# 清空 Wayland 相关变量
Environment="WAYLAND_DISPLAY="
Environment="WAYLAND_DEBUG=0"
Environment="QT_QPA_PLATFORM=eglfs"
Environment="GDK_BACKEND=x11"
Environment="SDL_VIDEODRIVER=offscreen"
Environment="XDG_SESSION_TYPE=tty"
Environment="LIBINPUT_DEFAULT_TOUCH_ENABLED=1"

# 访问 GPU/输入设备需要的附加组
SupplementaryGroups=input video render tty dialout

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
OOMScoreAdjust=-100

# 启动前的清理与检查
ExecStartPre=/bin/bash -c 'pkill -9 bamboo_integrated || true'
ExecStartPre=/bin/bash -c 'pkill -9 bamboo_cu || true'
ExecStartPre=/bin/bash -c 'for p in $(lsof -t /dev/video0 2>/dev/null); do kill -9 $p; done || true'
ExecStartPre=/bin/bash -c 'for p in $(lsof -t /dev/video1 2>/dev/null); do kill -9 $p; done || true'
ExecStartPre=/bin/bash -c 'systemctl restart nvargus-daemon'

# 停止并屏蔽 weston/nvweston/seatd，避免占用 /dev/dri
ExecStartPre=/bin/bash -c 'systemctl stop nvweston.service weston.service seatd.service 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl disable nvweston.service weston.service seatd.service 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'systemctl mask nvweston.service weston.service seatd.service 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'pkill -9 weston 2>/dev/null || true'
ExecStartPre=/bin/sleep 1

# 自定义停止逻辑
ExecStop=/bin/bash -c 'kill -TERM $MAINPID; sleep 2; pkill -9 -g $MAINPID bamboo_integrated || true'

# 重启与停止行为
KillMode=control-group
SendSIGKILL=yes
Restart=always
RestartSec=15
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=mixed

# 允许访问 GPU/设备
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# 输出到 journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
