# lvgl_frontend/scripts/bamboo-controller.service.in
[Unit]
Description=智能切竹机控制系统 - LVGL版本
Documentation=https://github.com/bambootech/bamboo-controller
After=network.target graphical.target
Wants=network.target
PartOf=graphical-session.target

[Service]
Type=simple
User=root
Group=root
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/bamboo_controller_lvgl -c @CMAKE_INSTALL_PREFIX@/config/default_config.json -f
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# 环境变量
Environment="LD_LIBRARY_PATH=/usr/local/cuda/lib64:/usr/lib/aarch64-linux-gnu"
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/0"

# 工作目录
WorkingDirectory=@CMAKE_INSTALL_PREFIX@

# 标准输出和错误输出
StandardOutput=append:/var/log/bamboo/controller.log
StandardError=append:/var/log/bamboo/controller_error.log

# 资源限制
LimitNOFILE=65536
LimitMEMLOCK=infinity

# 安全设置
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

# 设备访问权限
DeviceAllow=/dev/fb0 rw
DeviceAllow=/dev/input/event0 rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
DeviceAllow=/dev/nvidia-modeset rw
DeviceAllow=/dev/ttyUSB0 rw

[Install]
WantedBy=graphical.target
Alias=bamboo.service