# package/bamboo-cut-kiosk.service
[Unit]
Description=智能切竹机Kiosk模式界面
Documentation=https://github.com/bamboo-cut/bamboo-cut-develop
After=graphical-session.target gdm.service
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/bamboo-cut
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/usr/bin/bamboo-cut-kiosk
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# 确保在图形界面启动后运行
ExecStartPre=/bin/sleep 10

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bamboo-cut-kiosk

[Install]
WantedBy=graphical-session.target 