# package/bamboo-cut.service
[Unit]
Description=智能切竹机控制系统
Documentation=https://github.com/bamboo-cut/bamboo-cut-develop
After=network.target graphical-session.target
Wants=network.target

[Service]
Type=simple
User=bamboo-cut
Group=bamboo-cut
WorkingDirectory=/opt/bamboo-cut
Environment=PYTHONPATH=/opt/bamboo-cut
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/python3 /usr/bin/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bamboo-cut/logs /opt/bamboo-cut/temp /opt/bamboo-cut/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 资源限制
MemoryLimit=2G
CPUQuota=80%

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bamboo-cut

[Install]
WantedBy=multi-user.target 